# startup/crontab/crontab.entry
@reboot sleep 60 && ~/crontab/start_procServ_EPICS_RaspPi_IOC.sh